# deploy/video-to-essay-deliver.service
[Unit]
Description=Video-to-Essay Deliver Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/video-to-essay
EnvironmentFile=/root/video-to-essay/.env
ExecStart=/root/.local/bin/uv run python -c "from video_to_essay.deliver_worker import deliver_loop; deliver_loop()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
