# deploy/video-to-essay-download.service
[Unit]
Description=Video-to-Essay Download Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/video-to-essay
EnvironmentFile=/root/video-to-essay/.env
ExecStart=/root/.local/bin/uv run python -c "from video_to_essay.download_worker import download_loop; download_loop()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
